# ec2-management/cf-preserve.service
[Unit]
Description=Cloud Foundry VM State Preservation
Documentation=https://github.com/konveyor/asset-generation/blob/main/docs/deploying-cf-locally.md
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target
Requires=vboxdrv.service
After=vboxdrv.service
Conflicts=shutdown.target reboot.target halt.target
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStop=/usr/local/bin/cf-preserve.sh
TimeoutStopSec=1800
StandardOutput=journal
StandardError=journal

# Run as root for VirtualBox operations
User=root
Group=root

[Install]
WantedBy=multi-user.target
